# docker/supervisord-uwsgi.conf
[program:user_map]
command = /usr/local/bin/uwsgi --ini /home/web/user_map/docker/uwsgi.ini
user=web
stopsignal=INT

